# Nodejs/service/supervisor.conf
[program:site]
command=node .
directory=/home/njs/site
user=njs
